# debian/cm4all-ping-exporter.socket
[Unit]

[Socket]
ListenStream=/run/cm4all/prometheus-exporters/ping.socket

[Install]
WantedBy=sockets.target
