# debian/cm4all-fs-exporter.socket
[Unit]

[Socket]
ListenStream=%t/cm4all/prometheus-exporters/fs.socket

[Install]
WantedBy=sockets.target
