# debian/cm4all-cgroup-exporter.socket
[Unit]

[Socket]
ListenStream=/run/cm4all/prometheus-exporters/cgroup.socket

[Install]
WantedBy=sockets.target
