# debian/cm4all-process-exporter.socket
[Unit]

[Socket]
ListenStream=@cm4all-process-exporter

[Install]
WantedBy=sockets.target
