# debian/cm4all-multi-exporter.socket
[Unit]

[Socket]
ListenStream=5476
DeferAcceptSec=10
NoDelay=true

[Install]
WantedBy=sockets.target
