# debian/cm4all-process-exporter.socket
[Unit]

[Socket]
ListenStream=/run/cm4all/prometheus-exporters/process.socket

[Install]
WantedBy=sockets.target
