# debian/cm4all-cgroup-exporter.socket
[Unit]

[Socket]
ListenStream=@cm4all-cgroup-exporter

[Install]
WantedBy=sockets.target
