# debian/cm4all-kernel-exporter.socket
[Unit]

[Socket]
ListenStream=%t/cm4all/prometheus-exporters/kernel.socket

[Install]
WantedBy=sockets.target
